# plot-estate.timer
[Unit]
Description=Run plot-estate

[Timer]
OnCalendar=daily
Persistent=true   

[Install]
WantedBy=timer.target
